# systemd/anilist-discord.service
[Unit]
Description=Anilist Discord Bot Service
After=network.target

[Service]
Type=notify
User=$USER
Group=$USER
WorkingDirectory=/opt/anilist-randomizer-discord
Environment="NODE_ENV=production"
ExecStart=npm start
Restart=on-failure
RestartSec=10
TimeoutStopSec=60

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
